# socialhub.service
[Unit]
Description=SocialHub Server Service

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/local/go/bin/go run /home/ubuntu/github/socialhub/src/main.go
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target